# docker/supervisord/supervisord-prod.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/krakenhashes/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=20
priority=999
stdout_logfile=/var/log/krakenhashes/nginx/stdout.log
stderr_logfile=/var/log/krakenhashes/nginx/stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3
user=root

[program:backend]
command=/usr/local/bin/krakenhashes
directory=/etc/krakenhashes
autostart=true
autorestart=true
startsecs=10
startretries=5
priority=1
stdout_logfile=/var/log/krakenhashes/backend/stdout.log
stderr_logfile=/var/log/krakenhashes/backend/stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3
user=krakenhashes
environment=HOME="/home/krakenhashes",USER="krakenhashes"

[program:logrotate]
command=/bin/sh -c 'while true; do /usr/sbin/logrotate /etc/logrotate.d/krakenhashes >> /var/log/krakenhashes/logrotate.log 2>> /var/log/krakenhashes/logrotate.err; sleep 3600; done'
autostart=true
autorestart=true
stdout_logfile=/var/log/krakenhashes/logrotate.log
stderr_logfile=/var/log/krakenhashes/logrotate.err
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
stdout_logfile_backups=1
stderr_logfile_backups=1